# assets/systemd/zscaler-guardian-daemon.service
[Unit]
Description=start the zscaler-guardian daemon

[Service]
ExecStart=__EXECUTABLE_PATH__
Restart=always

[Install]
WantedBy=default.target
